# luci-app-ispapp/root/etc/init.d/ispapp
#!/bin/sh /etc/rc.common
## ispapp service manager ##

START=99
STOP=10
PID_FILE="/var/run/ispappd.pid"
LUA_SCRIPT="/bin/ispappd"

start() {
    echo "Starting ispappd..."
    
    # Start the Lua script in the background and save the PID
    /usr/bin/env lua $LUA_SCRIPT start &
    echo $! > $PID_FILE

    echo "ispappd started with PID $(cat $PID_FILE)"
}

stop() {
    if [ -f "$PID_FILE" ]; then
        PID=$(cat $PID_FILE)
        echo "Stopping ispappd with PID $PID..."
        
        # Kill the process and remove the PID file
        kill $PID && rm -f $PID_FILE
        echo "ispappd stopped."
    else
        echo "ispappd is not running or PID file not found."
    fi
}

restart() {
    echo "Restarting ispappd..."
    stop
    sleep 1
    start
}

reload() {
    if [ -f "$PID_FILE" ]; then
        PID=$(cat $PID_FILE)
        echo "Reloading ispappd with PID $PID..."
        
        # Send SIGHUP to the process to reload configuration
        kill -HUP $PID
        echo "ispappd reloaded."
    else
        echo "ispappd is not running or PID file not found."
    fi
}
