# luci-app-ispapp/root/etc/uci-defaults/50_luci-ispapp
#!/bin/sh
# reload luci cache
rm -rf /tmp/luci-*
exit 0
return 0
