# luci-app-ispapp/root/etc/uci-defaults/50_luci-ispapp
#!/bin/sh
# reload luci cache
rm -rf /tmp/luci-*
if [ -x /etc/init.d/uhttpd ]; then
/etc/init.d/uhttpd restart
fi
if [ -x /etc/init.d/ispapp ]; then
	chmod +x /etc/init.d/ispapp
	/etc/init.d/ispapp enable
	/etc/init.d/ispapp start
fi
if [ -f /bin/ispappd ]; then
	chmod +x /bin/ispappd
fi
if [ -f /usr/libexec/rpcd/ispapp ]; then
	chmod +x /usr/libexec/rpcd/ispapp
	/etc/init.d/rpcd reload
fi
exit 0
return 0
