# luci-app-ispapp/root/etc/uci-defaults/50_luci-ispapp
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ispapp[-1]
	add ucitrack ispapp
	set ucitrack.@ispapp[-1].init=ispapp
	commit ucitrack
EOF

chmod +x /etc/init.d/ispapp /usr/libexec/rpcd/ispapp /bin/ispappd

rm -f /tmp/luci-indexcache
exit 0